# scripts/toune-library-sync.service
[Unit]
Description=Toune-o-matic Library Links Sync
After=network.target toune-backend.service
Wants=toune-backend.service

[Service]
Type=oneshot
User=pi
Group=pi
Environment="TOUNE_API_BASE=http://127.0.0.1:11000/api"
ExecStart=/srv/toune/repo/toune-o-matic/scripts/sync-library-links.sh
